# scripts/certbot-renew.cron
31 23 * * * root docker compose run --rm certbot renew --deploy-hook "docker exec nginx nginx -s reload" >> /var/log/cron-certbot.log 2>&1
